# containers/rootfs/etc/s6-overlay/s6-rc.d/cron/run
#!/command/execlineb -P
with-contenv
s6-notifyoncheck -d
/usr/sbin/cron -f -l -L 15